# services/i3-session.target
[Unit]
Description=i3 session running
BindsTo=graphical-session.target